# peercoin-seeder.service
[Unit]
Description=Peercoin DNS Seeder
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/peercoin-seeder
ExecStart=/opt/peercoin-seeder/peercoin-seeder-rust
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/peercoin-seeder

# Network permissions for DNS port 53
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
